# deploy/systemd/github-runner.service
[Unit]
Description=GitHub Actions Runner
After=network.target

[Service]
ExecStart=/home/samuel/actions-runner/run.sh
WorkingDirectory=/home/samuel/actions-runner
Restart=no
User=samuel

[Install]
WantedBy=multi-user.target