# deploy/systemd/oven-kiosk.service
[Unit]
Description=Oven UI Kiosk
After=graphical.target network-online.target

[Service]
User=samuel  # Change this if necessary
Environment=DISPLAY=:0  # Only needed for X11
Environment=XDG_RUNTIME_DIR=/run/user/1001  # Only needed for Wayland
ExecStart=/usr/bin/chromium-browser \
  --noerrdialogs --kiosk --incognito --start-fullscreen \
  --disable-restore-session-state --disable-infobars \
  http://localhost:3000
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
