# deploy/systemd/oven-kiosk.service
[Unit]
Description=Oven UI Kiosk
Requires=graphical.target
After=graphical.target

[Service]
User=samuel
ExecStart=/home/smart-oven/oven-kiosk.sh
Restart=no

[Install]
WantedBy=graphical.target
