# deploy/systemd/oven-kiosk.service
[Unit]
Description=Oven UI Kiosk
After=graphical-session.target

[Service]
ExecStart=/usr/bin/chromium-browser \
  --kiosk --incognito --start-fullscreen --noerrdialogs \
  --disable-restore-session-state \
  http://localhost:3000
Restart=no

[Install]
WantedBy=default.target